# squid.cron
# Run the squid job every night
* * * * * root /usr/sbin/logrotate -f /squid.logrotate
