# squid.cron
# Run the squid job every night
59 23 * * * root /usr/sbin/logrotate -f /squid.logrotate
